# mysql-backup/crontab
0 20 * * 1-5 (. /lib.sh; dbparam; backupAll /mnt)
